# systemd/pi-client.service
[Unit]
Description=Pi Heartbeat Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mypi
WorkingDirectory=/home/mypi/EC_TG_Bot_v2
ExecStart=/bin/bash -lc "cd /home/mypi/EC_TG_Bot_v2 && source .venv/bin/activate && python3 entrypoints/start_pi_client.py"
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
